# compose-client/echo/cron.d/cron_echo.sh
#!/bin/sh

export VERSION_ENV="version1"

echo VERSION_ENV: ${VERSION_ENV}
