# compose-client/echo/cron.d/echo_version
* * * * * dev-user /bin/bash /etc/cron.d/cron_echo.sh >> /tmp/cron_echo/crontest.txt
